# gen-config/supervisor-devpi.conf
[program:devpi-server]
command=F:\Projects\Python\package_t\.venv\Scripts\devpi-server 
priority = 999
startsecs = 5
redirect_stderr = True
autostart = True
